# 89-supervisord-mtm.ini

[program:mtm]
directory=%(here)s/
command=make -C %(here)s/../src/browenv open-mtm
autostart=true
autorestart=true
