# 90-supervisord-brw.ini

[program:brw]
directory=%(here)s/
command=make -C %(here)s/../src/browenv serve-brw-pgs
autostart=true
autorestart=true
