# 91-supervisord-pbm.ini

[program:pbm]
directory=%(here)s/
command=make -C %(here)s/../src/browenv serve-pbmweb
autostart=true
autorestart=true
